# root/etc/services.d/web/finish
#!/usr/bin/with-contenv sh

s6-svc -d /var/run/s6/services/web